# debian/bunsen-apt-update-checker.user.service
[Unit]
Description=Check for available apt upgrades

[Service]
Type=oneshot
ExecStart=bl-apt-update-check
